# internal/setup/templates/units/launchd.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cc-handoff.watch</string>

  <key>ProgramArguments</key>
  <array>
    <string>{{.BinPath}}</string>
    <string>watch</string>
  </array>

  <key>WorkingDirectory</key>
  <string>{{.WorkDir}}</string>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/cc-handoff.watch.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/cc-handoff.watch.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
